# dvwa/supervisord.conf
[supervisord]
nodaemon=true

[program:dvwa]
command=apache2-foreground
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
